# meta-image-common/recipes-core/mount-data/files/cubieboard2/data.mount
[Unit]
Description=Mount /data
ConditionPathIsDirectory=/data
Before=local-fs.target
After=swap.target

[Mount]
What=/dev/mmcblk0p4
Where=/data
Type=ext4
Options=defaults
